# .claude/schedule/com.claude.weekly-worksheet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude.weekly-worksheet</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/stefania/Spec Driven Development/hours-registration/.claude/schedule/run-weekly-worksheet.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>5</integer>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/stefania/Spec Driven Development/hours-registration</string>
  <key>StandardOutPath</key>
  <string>/Users/stefania/.claude/jobs/weekly-worksheet.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/stefania/.claude/jobs/weekly-worksheet-error.log</string>
</dict>
</plist>
